# Update to 1.6.1/static_files/etc/supervisord.d/extra.ini
[program:mlmmjadmin]
startsecs = 0
autorestart = true
command=/services/mlmmjadmin.sh

[program:memcached]
startsecs = 0
autorestart = true
command=/services/memcached.sh

[group:extra]
programs=mlmmjadmin,memcached
priority=800