# mysql/static_files/etc/supervisord.d/extra.ini
[program:mlmmjadmin]
startsecs = 0
autorestart = true
command=/services/mlmmjadmin.sh


[group:extra]
programs=mlmmjadmin
priority=800