# mysql/static_files/etc/supervisord.d/mariadb-server.ini
[program:mariadb-server]
startsecs = 0
autorestart = true
command=/services/mariadb-server.sh


# TODO: group restart seems to be problematic
[group:db-server]
programs=mariadb-server
priority=100