# mysql/static_files/etc/supervisord.d/mariadb-server.ini
[program:mariadb-server]
startsecs = 0
autorestart = true
command=/etc/services/mariadb-server.sh