# mysql/static_files/etc/supervisord.d/mta.ini
[program:postfix]
startsecs = 0
autorestart = true
command=/services/postfix.sh
priority=210

[program:clamav-daemon]
startsecs = 0
autorestart = true
command=/services/clamav-daemon.sh
priority=220

[program:amavis]
startsecs = 0
autorestart = true
command=/services/amavis.sh
priority=230

[program:dovecot]
startsecs = 0
autorestart = true
command=/services/dovecot.sh
priority=240


# GROUP definition
[group:mta]
programs=postfix,amavis,clamav-daemon,dovecot
priority=200